# goodies/fruitypipy.conf
[program:fruitypipy]
directory=/opt/FruityPiPy
command=python main.py
autostart=true
autorestart=true
